# plugins/telegram-proxy/telegram-proxy.service.erb
# Put this file in /usr/lib/systemd/system

[Unit]
Description=Telegram MQTT+HTTP Proxy
Requires=network.target

[Service]
User=david
Environment=TELEGRAM_BOT_TOKEN=<%= ENV['TELEGRAM_BOT_TOKEN'] %>
Environment=TELEGRAM_CHAT_ID=<%= ENV['TELEGRAM_CHAT_ID'] %>
ExecStart=<%= ENV['NODE'] %> <%= ENV['HOME'] %>/.dotfiles/scripts/telegram-proxy
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure

[Install]
WantedBy=multi-user.target
