# services/packages.service.erb
# Put this file in /usr/lib/systemd/user

[Unit]
Description=Packages server

[Service]
User=david

ExecStart=<%= ENV['HOME'] %>/.dotfiles/scripts/pkgs server
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure

[Install]
WantedBy=multi-user.target
